# betalightboard.service
[Unit]
Description=Beta Light Board
# Wait for real network so the start script can fetch dependency updates on boot.
After=network-online.target
Wants=network-online.target

[Service]
User=root
WorkingDirectory=/home/pi/BetaLightBoard-pi
# Entrypoint lives in the repo (scripts/start.sh) so the boot sequence —
# dependency install + launch — updates via git/the update button without
# ever needing to re-edit this unit. Allow time for pip on a slow first boot.
ExecStart=/home/pi/BetaLightBoard-pi/scripts/start.sh
TimeoutStartSec=300

Restart=always
RestartSec=5
Environment=LED_SIMULATE=0

[Install]
WantedBy=multi-user.target
